# bluetooth-audio/udev-rules/99-bluetooth-udev.rules
SUBSYSTEM=="input", KERNEL=="input[0-9]*", RUN+="/usr/src/bluetooth-udev"
